# src/Dotnet.Installer.Console/Scripts/{SNAP}-update-watcher.path
[Unit]
Description=dotnet-installer-update-mounts watcher for {SNAP}

[Path]
PathExists=/var/snap/{SNAP}/common/dotnet-installer-update-mounts

[Install]
WantedBy=multi-user.target
