# src/Dotnet.Installer.Console/Scripts/{SNAP}-update-watcher.path
[Unit]
Description=needs-update watcher for {SNAP}

[Path]
PathExists=/var/snap/{SNAP}/common/needs-update

[Install]
WantedBy=multi-user.target
